# root/etc/services.d/traccar/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid dockeruser java -jar $TRACCAR_VOL/tracker-server.jar $TRACCAR_VOL/conf/traccar.xml
